# automation/com.beibei.daily-site-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.beibei.daily-site-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/apple/.cache/codex-runtimes/codex-primary-runtime/dependencies/python/bin/python3</string>
    <string>/Users/apple/Documents/贝贝外刊/scripts/build_beibei_daily_site.py</string>
    <string>--watch</string>
    <string>--interval</string>
    <string>30</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/apple/Documents/贝贝外刊</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/apple/Documents/贝贝外刊/.beibei-cache/automation.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/apple/Documents/贝贝外刊/.beibei-cache/automation-error.log</string>
</dict>
</plist>
